# supervisor/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx

[program:uwsgi]
command =/opt/conda/bin/uwsgi --ini  /var/www/app/uwsgi.ini
